# service/log/run
#!/bin/sh
exec envuidgid mersenne multilog /var/log/mersenne
