# service/log/run
#!/bin/sh
exec envuidgid mersenne multilog t s10485760 '!tai64nlocal' /var/log/mersenne
